# scripts/storebackend.service
[Unit]
Description=Store Backend Spring Boot Application
After=network.target postgresql.service
Requires=postgresql.service

[Service]
Type=simple
User=storebackend
Group=storebackend
WorkingDirectory=/opt/storebackend

# Load all app config + secrets from this file (written by deploy.sh)
EnvironmentFile=/etc/storebackend.env

# Start the Spring Boot JAR
ExecStart=/usr/bin/java $JAVA_OPTS -jar /opt/storebackend/app.jar

# Restart on failure
Restart=on-failure
RestartSec=10

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=storebackend

# Limits
LimitNOFILE=65536
LimitNPROC=4096

[Install]
WantedBy=multi-user.target
