# scripts/storebackend.service
[Unit]
Description=Store Backend Spring Boot Application
# Warten auf Netzwerk (online, nicht nur "up"), PostgreSQL und MinIO
After=network-online.target postgresql.service minio.service
Wants=network-online.target
Requires=postgresql.service

[Service]
Type=simple
User=storebackend
Group=storebackend
WorkingDirectory=/opt/storebackend

# Alle App-Konfigurationen + Secrets aus dieser Datei (von deploy.sh/setup-vps.sh geschrieben)
EnvironmentFile=/etc/storebackend.env

# $JAVA_OPTS kommt aus EnvironmentFile (inkl. Xms/Xmx/HeapDump-Flags)
ExecStart=/usr/bin/java $JAVA_OPTS -jar /opt/storebackend/app.jar

# Neustart bei Fehler; 20s Wartezeit damit DB/MinIO nach Reboot stabil sind
Restart=on-failure
RestartSec=20

# Logging via systemd journal
StandardOutput=journal
StandardError=journal
SyslogIdentifier=storebackend

# systemd Memory-Limits
# MemoryHigh: Ab hier beginnt Kernel aggressiv zu GC-drücken (sanft)
# MemoryMax:  Hartes Limit inkl. Java Native Memory, Metaspace, Threads, DirectBuffers.
#              3G für Heapdump-Diagnose (Java Xmx=1024m + ~1G Native + Puffer).
#              Nach OOM-Analyse auf 2G reduzieren.
MemoryHigh=1600M
MemoryMax=3G

# Datei-/Prozess-Limits
LimitNOFILE=65536
LimitNPROC=4096

[Install]
WantedBy=multi-user.target
